# deploy/bharatquant-dashboard.service
[Unit]
Description=BharatQuant FastAPI dashboard
After=network.target bharatquant-engine.service

[Service]
Type=simple
User=bharatquant
WorkingDirectory=/opt/bharatquant/zerodha-momo-rl
EnvironmentFile=/etc/bharatquant/env
ExecStart=/usr/bin/python3.11 -m src.api.dashboard
Restart=on-failure

[Install]
WantedBy=multi-user.target
